# Library/LaunchDaemons/root/locallan.allowed-http.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>locallan.allowed-http.update</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>/usr/local/bin/parallel "/usr/bin/host {} | /usr/bin/awk '\$0 ~ /has.*address/ {print \$NF}'" :::: /usr/local/etc/allowed-http.names | pfctl -a locallan.home -T replace -t allowed-http -f -</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>3600</integer>
</dict>
</plist>

